# clock/prayer_times.crontab
59 13 9 3 * catt -d "home speaker" cast "https://ia800304.us.archive.org/25/items/adhan_202002/adhan.mp3" > /dev/null 2>&1 # test_prayer_time_job
