# clock/prayer_times.crontab
53 14 9 3 * catt -d "home speaker" cast "https://www.youtube.com/watch?v=2-XRbcLQ6b8" > /dev/null 2>&1 # test
